# com.claude.usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.usage</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the full path returned by `which uv` -->
        <string>/Users/YOUR_USERNAME/.local/bin/uv</string>
        <string>run</string>
        <string>/Users/YOUR_USERNAME/code/ai/claude_usage_mac_menu/main.py</string>
    </array>

    <!-- The API key is stored in the macOS Keychain (set via Settings → Set API Key…).
         You only need this section if you prefer to pass the key as an env var instead. -->
    <!-- <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_ADMIN_KEY</key>
        <string>sk-ant-admin-REPLACE-ME</string>
    </dict> -->

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/claude_usage.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/claude_usage.err</string>
</dict>
</plist>
